# leque/leque.service
[Unit]
Description=AMDGPU fan controller

[Service]
Type=exec
Restart=always
ExecStart=/usr/bin/leque

[Install]
WantedBy=multi-user.target
